# crontab
# m h  dom mon dow   command
0 * * * * cd /path/to/sync_family && wget -N https://example.com/path/to/family

